# check_mk@.service
[Unit]
Description=Check_MK LibreNMS Agent Service
After=network.target
 
[Service]
ExecStart=/usr/bin/check_mk_agent
StandardOutput=socket
