# check_mk.socket
[Unit]
Description=Check_MK LibreNMS Agent Socket
 
[Socket]
ListenStream=6556
Accept=yes
 
[Install]
WantedBy=sockets.target
